# files/supervisor/sophon.conf
[program:sophon]
command=/app/sophon/venv/bin/sophon
directory=/app/sophon
stdout_logfile=/app/sophon/sophon.log
stderr_logfile=/app/sophon/sophon.err
autostart=true
autorestart=true
